# Build/dnscrypt1.6.1_build_docker/dnscrypt.service
[Unit]
Description=dnscrypt
After=network.target iptables.service firewalld.service

[Service]
ExecStart=/bin/dnscrypt-start

[Install]
WantedBy=multi-user.target
